# frstack/rootfs/etc/services.d/frstack/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start FRStack
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Start FRstack
bashio::log.info "Starting FRStack"
exec dotnet /usr/local/frstack/FRStackWebApi.dll \
  --urls "*:5025"
