# frstack/rootfs/etc/services.d/frstack/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start FRStack
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Start FRstack
bashio::log.info "Starting FRStack version ${API_VERSION}"
bashio::log.debug "Arch is ${BUILD_ARCH}"
bashio::log.debug "FRStack port is ${PORT}"

#exec /usr/local/bin/frstackwebapi
